# mcloud/upstart.conf

description "Mcloud server"
author "Modera"

start on filesystem
respawn
script
  /opt/mcloud/bin/mcloud-server > /var/log/mcloud.log 2>&1
end script
